# supervisor/tailscaled.conf
[program:tailscaled]
command=/app/tailscaled
directory=/app
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/tailscaled/tailscaled.log